# ansible/roles/glances/files/glances.service
[Unit]
Description=Glances
After=network.target

[Service]
User=glances
ExecStart=
Restart=on-abort

[Install]
WantedBy=multi-user.target
