# ansible/roles/init_wifi/templates/dhclient.service.jinja2
[Unit]
Description= DHCP Client
Before=network.target

[Service]
Type=simple
ExecStart=/sbin/dhclient {{ wifi.interface }}

[Install]
WantedBy=multi-user.target
