# ansible/roles/homeassistant_init/files/home-assistant@homeassistant.service
[Unit]
Description=Home Assistant
After=network-online.target

[Service]
Type=simple
User=%i
ExecStart=

[Install]
WantedBy=multi-user.target
